# startup/dot-config/systemd/user/startup.service
[Unit]
Description=Startup Script

[Service]
ExecStart=%h/.startup.sh

[Install]
WantedBy=default.target
